# supervisord.conf
[supervisord]
nodaemon=true

[program:one]
command=echo 'test'

[program:two]
command=echo 'test 2'
